# etc/init.d/icd2
#!/sbin/openrc-run
# See the git history of this file for the old (shell-script) version.

description="Internet Connectivity daemon 2"

depend() {
	need networking dsme
}

start_pre() {
	ebegin "Starting $description"
	/usr/sbin/waitdbus system
	. /tmp/session_bus_address.user
	/usr/sbin/waitdbus session
	/usr/bin/killall -9 udhcpc >/dev/null 2>&1 || true
}

start() {
	/usr/sbin/dsmetool -m -17 -f "/usr/sbin/icd2 -l2"
}

stop() {
	eend "Stopping $description"
	/usr/sbin/dsmetool -k "/usr/sbin/icd2 -l2"
}
