# debian/postinst
#!/bin/sh

rc-update add icd2 default || true

exit 0
