# debian/prerm
#!/bin/sh

rc-update del icd2 default || true

exit 0
